# build/wr_runtime/pkg/etc/init.d/whiterails
#!/sbin/openrc-run
description="WhiteRails semantic runtime"
command="/usr/bin/wr_runtime"
command_background="yes"
pidfile="/run/wr_runtime.pid"

depend() {
    need localmount
}
